# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

; ── TigerVNC server ────────────────────────────────────────────────────────────
[program:tigervnc]
command=/entrypoint.sh vnc
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/tigervnc.log
stderr_logfile=/var/log/supervisor/tigervnc.err
environment=HOME="/home/rosdev",USER="rosdev"

; ── noVNC websocket proxy ──────────────────────────────────────────────────────
[program:novnc]
command=/usr/bin/websockify --web /usr/share/novnc/ %(ENV_NOVNC_PORT)s localhost:%(ENV_VNC_PORT)s
autostart=true
autorestart=true
priority=20
startsecs=3
stdout_logfile=/var/log/supervisor/novnc.log
stderr_logfile=/var/log/supervisor/novnc.err
